# 90-honor-fmi721-kbd.rules
SUBSYSTEM=="leds", KERNEL=="honor::kbd_backlight", RUN+="/bin/chmod 0666 /sys/class/leds/%k/brightness"
